# semispace-comet/semispace-comet-app/src/main/resources/init.d/jetty.gentoo
#!/sbin/runscript

depend() {
  need net
}

# The following configuration should reside in /etc/conf.d/...
# Please move it to the relevant file.
# export JETTY_HOME=/opt/semispace/app
# export JETTY_RUN=${JETTY_HOME}
# export TMP=${JETTY_HOME}/${TMP}
# export JAVA_OPTIONS="-server -Xms256m -Xmx368m -XX:PermSize=128m -XX:MaxPermSize=256m -XX:NewSize=64m -XX:MaxNewSize=128m"
# export JETTY_USER="semispace"
# export JAVA_HOME=`java-config -o`


start() {
  ebegin "Starting semispace"
  if [ -z "${JETTY_USER}" ]; then
      eerror "Need to have the JETTY_USER variable set"
      return 1
  fi
  if [ -z "${JETTY_HOME}" ]; then
      eerror "Need to have the JETTY_HOME variable set"
      return 1
  fi

  ${JETTY_HOME}/bin/jetty.sh start
  eend $?
}

stop() {
  ebegin "Stopping semispace"
    ${JETTY_HOME}/bin/jetty.sh stop
  eend $?
}
